# scripts/suns.service
# A systemd service file to run the Hugo development server for the SUNS website.
#
# Install with
# ln -sf $HOME/work/SUNS/scripts/suns.service $HOME/.config/systemd/user/suns.service
# Then enable and start with
# systemctl --user enable --now suns.service

[Unit]
Description=SUNS Hugo Development Server
After=network.target

[Service]
WorkingDirectory=%h/work/SUNS/www
Environment=PATH=%h/.bun/bin:%h/.local/bin:/usr/local/bin:/usr/bin
Restart=always
ExecStart=%h/.local/bin/hugo server --bind 0.0.0.0 --port 7867 --baseURL https://suns.chineseroom.micahrl.com/ --appendPort=false --environment development

[Install]
WantedBy=default.target

[Install]
WantedBy=default.target
